# systemd/cogitator-serial.service
[Unit]
Description=Calvin Cogitator Serial Service
After=cogitator-broker.service
Requires=cogitator-broker.service

[Service]
Type=simple
WorkingDirectory=/opt/cogitator
Environment=PYTHONPATH=/opt/cogitator
ExecStart=/usr/bin/python3 services/serial/serial_service.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
